# files/zabbix-backup.timer
[Unit]
Description=Backup postgresql data of zabbix

[Timer]
OnCalendar=*-*-* 3,15:00
Persistent=true

[Install]
WantedBy=timers.target
