# files/zabbix-backup.service
[Unit]
Description=Backup postgresql data of zabbix

[Service]
Type=oneshot
ExecStart=/etc/zabbix/backup.sh
